# nginx/test/supervisord.conf
# ─────────────────────────────────────────────────────────────────────────────
# supervisord.conf  — agentic_nginx_prod_test
#
# supervisord runs as PID 1.  Each child can be killed and restarted
# independently — the container itself stays alive.
#
# Process start order (priority):
#   10  postgresql  — DB must be ready before testapp connects
#   15  sshd        — SSH must be up before watcher tries host_service_* actions
#   20  nginx       — web server (depends on testapp being up for /health proxy)
#   30  testapp     — application backend
# ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

# Allow any local user (incl. ops via systemctl shim) to run supervisorctl.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0777

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── PostgreSQL ────────────────────────────────────────────────────────────────
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/data -c logging_collector=on -c log_directory=/var/log/postgresql -c log_filename=postgresql.log
user=postgres
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/supervisor/postgresql-stdout.log
stderr_logfile=/var/log/supervisor/postgresql-stderr.log
priority=10
stopsignal=INT
stopwaitsecs=30
startretries=3

# ── SSH daemon ────────────────────────────────────────────────────────────────
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd-err.log
priority=15
stopsignal=TERM

# ── nginx ─────────────────────────────────────────────────────────────────────
# 'daemon off' keeps nginx in the foreground so supervisord tracks it.
# nginx is NOT at PID 1 — killing it (pkill nginx, systemctl restart nginx,
# SIGKILL from a runbook) just restarts this process; the container lives on.
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx-stdout.log
stderr_logfile=/var/log/supervisor/nginx-stderr.log
priority=20
stopsignal=QUIT
stopwaitsecs=10
startretries=5

# ── testapp ───────────────────────────────────────────────────────────────────
[program:testapp]
command=/usr/bin/python3 /opt/testapp/app.py
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/app/testapp.log
stderr_logfile=/var/log/app/testapp-err.log
priority=30
stopsignal=TERM
startretries=5
